# config/test/supervisor/conf.d/celery-worker.conf
[program:tz-data-dws-celery-worker]
directory = %(ENV_DJANGO_APP_HOME)s/mysite ; 程序的启动目录
command = /usr/bin/celery -A mysite worker --loglevel=DEBUG ; 启动命令

redirect_stderr = true  ; 把 stderr 重定向到 stdout，默认 false
stdout_logfile = %(ENV_DJANGO_APP_HOME)s/logs/celery-worker.log
