# config/test/supervisor/conf.d/celery-beat.conf
[program:celery-beat]
directory = %(ENV_DJANGO_APP_HOME)s/mysite ; 程序的启动目录
command = /usr/bin/celery -A mysite beat --pidfile %(ENV_DJANGO_APP_HOME)s/running/handle/celerybeat.pid --schedule %(ENV_DJANGO_APP_HOME)s/running/handle/celerybeat-schedule --loglevel=DEBUG ; 启动命令

redirect_stderr = true  ; 把 stderr 重定向到 stdout，默认 false
stdout_logfile = %(ENV_DJANGO_APP_HOME)s/logs/celery-beat.log
